# conf/supervisor_users_online.conf
[supervisord]
minprocs=3000

[program:users_online]
command=php scripts/create_user_session.php
numprocs=2000
process_name=user_session%(process_num)s
autostart=true
autorestart=true
startsecs=0
